# config/supervisor/conf.d/flask_demo_flower.conf
; ================================
; celery flower supervisor example
; ================================

[program:flask_demo_flower]
environment=FLASK_CONFIG=%(ENV_FLASK_CONFIG)s
; Set full path to celery program if using virtualenv
command=celery flower -A run_celery.celery -l info --address=0.0.0.0 --port=5555

; remove the -A myapp argument if you are not using an app instance

directory=/var/www/flask_demo/
user=root
numprocs=1
stdout_logfile=/var/www/flask_demo/app/log/flask_demo_flower.log
stderr_logfile=/var/www/flask_demo/app/log/flask_demo_flower_error.log
autostart=true
autorestart=true
startsecs=10

; if rabbitmq is supervised, set its priority higher
; so it starts first